# tetrisd/tetrisd.service
[Unit]
Description=Tetris Daemon

[Service]
Type=forking
ExecStart=/usr/local/bin/tetrisd -d

[Install]
WantedBy=multi-user.target
